# build_scripts/server/debian_packaging/py_support/debian/postinst
#! /bin/sh

update-python-modules  openmolar-server.public
update-rc.d openmolar defaults

service openmolar restart
