# build_scripts/server/debian_packaging/py_support/debian/postinst
#! /bin/sh

#DEBHELPER#

update-python-modules  openmolar-server.public

openmolar-alter-master-user

update-rc.d openmolar defaults

service openmolar restart
