# build_scripts/server/debian_packaging/dh_python2/debian/postinst
#! /bin/sh

update-rc.d openmolar defaults

/etc/init.d/openmolar restart
